# .jobs.cron
0 *	* *	* if [ -f core.* ]; then rm core.*; fi
0 */2 *	* *	cd ~/.f; git pull origin master
