# .jobs.cron
*/5	*	*	*	*	for f in ~/core.*; do rm "$f" >& /dev/null; done
0	*/2	*	* 	*	git pull origin master
