# .jobs.cron
*/2	*	* 	*	*	cd ~/.f; git pull origin master >& /dev/null
